# service/ipfs.service
[Unit]
Description=Deamon ipfs

[Service]
Type=simple
User=lordpax
Group=wheel
ExecStart=/usr/bin/ipfs daemon

[Install]
WantedBy=multi-user.target
